# idlewatcher.service
[Unit]
Description=Idlewatcher

[Service]
ExecStart=/path/to/idlewatcher -t 300

[Install]
WantedBy=default.target
